# vagrant/files/usr/lib/systemd/system/kafka-connector.service
[Unit]
Description=Confluent Kafka Connector Distributed services
After=syslog.target network.target zookeeper.target kafka-rest.target

[Service]
LimitMEMLOCK=infinity
LimitNOFILE=65535

#Type=simple
# idle: similar to simple, the actual execution of the service binary is delayed
#       until all jobs are finished, which avoids mixing the status output with shell output of services.
RemainAfterExit=yes
User=root
Group=root
ExecStart=/etc/init.d/kafka-control start-kafka-connect
ExecStop=/etc/init.d/kafka-control stop

[Install]
WantedBy=multi-user.target
